# apache/files/pag
#! /usr/bin/pagsh
### BEGIN INIT INFO
# Provides:             pag
# Required-Start:       openafs-client
# Required-Stop:        openafs-client
# Should-Start:         $syslog $named openafs-fileserver
# Should-Stop:          openafs-fileserver
# X-Start-Before:		openafs-client
# X-Stop-After:         autofs
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    pag
# Description:          Starts or stops k5start for service keytab, 
#						getting the service tokens as necessary. 
#						Outside PAG only
### END INIT INFO
#
# Modified by Frederiko Costa<fcosta@stanford.edu> for Debian
# Copyright 2000, International Business Machines Corporation and others.
# All Rights Reserved.
# 
# This wrapper script will create a new pag and start both k5start and apache. Nothing else !!!

get_svc_name() {
        if [ -r /etc/redhat-release ]; then
                echo "httpd"
        elif [ -f /etc/lsb-release ] &&  grep -qi ubuntu /etc/lsb-release; then
                echo "apache2"
        else
				# defaults to Ubuntu
              	echo "apache2"
        fi
}

# define the service name to run
svcname=`get_svc_name`

case "$1" in
	start)
		if ! (/etc/init.d/$svcname status > /dev/null) && ! (/etc/init.d/k5start status >/dev/null); then
			/etc/init.d/k5start start 
			/etc/init.d/$svcname start		
		elif ! /etc/init.d/$svcname status > /dev/null || ! /etc/init.d/k5start status >/dev/null; then
			/etc/init.d/k5start restart
			/etc/init.d/$svcname restart
		fi
		;;
	stop)
		/etc/init.d/$svcname stop
		/etc/init.d/k5start stop
		;;
	status)
		/etc/init.d/$svcname status
		retval=$?
        [ $retval -ne 0 ] && exit $retval
		/etc/init.d/k5start status
		retval=$?
        [ $retval -ne 0 ] && exit $retval
		exit 0
		;;
	restart)
		/etc/init.d/k5start restart
		/etc/init.d/$svcname restart
		;;
	*)
		echo "Usage: $0 {start|stop|status|restart}"
		exit 2
		;;
esac
